# third-party/avr/lib/systemd/system/geoclue.service
[Unit]
Description=Location Lookup Service

[Service]
Type=dbus
BusName=org.freedesktop.GeoClue2
User=geoclue
ExecStart=/usr/lib/geoclue-2.0/geoclue -t 5
